# configs/gunicorn.conf
[program:gunicorn]
command=gunicorn app:app -b localhost:5000 --access-logfile /deploy/app/logs/access.log --error-logfile /deploy/app/logs/error.log
directory=/deploy/app
